# installer/caldera.service.d/xdr-auth-debug.conf
[Service]
Environment=XDR_CALDERA_AUTH_DEBUG=1
