# installer/xdr-lab-web-console@.service
[Unit]
Description=XDR Lab web console for %i
Documentation=file:///opt/xdr-lab/docs/web-console.md
After=network-online.target libvirtd.service virtqemud.service xdr-lab-host-network.service
Wants=network-online.target libvirtd.service xdr-lab-host-network.service
ConditionPathExists=/opt/xdr-lab/scripts/xdr-lab-vm-manager.sh

[Service]
Type=simple
Environment=XDR_LAB_WEB_CONSOLE_BIND=127.0.0.1
Environment=XDR_LAB_WEB_CONSOLE_PORT=6080
ExecStart=/opt/xdr-lab/scripts/xdr-lab-vm-manager.sh web-console serve %i
Restart=always
RestartSec=5
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
